# deb/drone/etc/init/drone.conf
start on (filesystem and net-device-up)

chdir /var/lib/drone
console log

script
    droned --port=":80"
end script