# dist/drone/etc/init/drone.conf
start on (filesystem and net-device-up)

chdir /var/lib/drone
console log

script
    drone --config /etc/drone/drone.toml
end script