# dist/drone/etc/init/drone.conf
start on (filesystem and net-device-up)

chdir /var/lib/drone
console log

script
	set -a
	if [ -f /etc/drone/dronerc ]; then
		. /etc/drone/dronerc
	fi
	/usr/local/bin/drone
end script
