# packaging/root/usr/share/drone/init/drone.conf
start on (filesystem and net-device-up)

chdir /var/lib/drone
console log

script
	/usr/local/bin/droned --config=/etc/drone/drone.toml
end script
